# node.ini
[program:publishing]
command=node /var/publishing/index.js
autostart=true
autorestart=true
stderr_logfile=/var/log/error-publishing.log
stdout_logfile=/var/log/out-publishing.log
